# build-rootfs/files/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

hvc0::respawn:/sbin/getty -L -l /usr/local/bin/podroid-login 115200 hvc0 xterm-256color
::shutdown:/sbin/openrc shutdown
